# app/jobs/crontab
# START CRON JOB
0 */1 * * * /usr/local/bin/python /app/send_email.py
* * * * * /usr/local/bin/python /app/add_file.py
0 12 * * * /usr/local/bin/python /app/suppression.py
# END CRON JOB